# root/var/spool/cron/crontabs/goldmund.crontab
0 20 * * 6 rm -rf ~/.cache/eww_covers/*
